# systemd/ainap.service
[Unit]
Description=AINAP - AI News Auto Publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ainap
Group=ainap
WorkingDirectory=/opt/ainap
ExecStart=/opt/ainap/venv/bin/python -m src.main
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/ainap/.env

# リソース制限（512MB保護）
MemoryMax=380M
MemoryHigh=300M
CPUQuota=80%

# セキュリティ
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/ainap/data /opt/ainap/logs

[Install]
WantedBy=multi-user.target
